# scripts/com.vcc.briefing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vcc.briefing</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>/ABSOLUTE/PATH/TO/vitals/scripts/brief-if-ready.sh</string>
    </array>

    <!-- EXAMPLE (macOS). Replace /ABSOLUTE/PATH/TO/vitals above with your repo
         path before `launchctl load`. Fires every 30 min from 07:00 to 12:00
         local; the script is idempotent (first window with data for today
         writes the briefing, later windows short-circuit). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/vitals/data/briefings/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/vitals/data/briefings/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
